# ansible/roles/lxc_node/templates/ebpf-blockchain.service.j2
[Unit]
Description=eBPF Blockchain Node
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory={{ node_working_dir }}
ExecStart={{ node_binary_dir }}/ebpf-node --iface {{ node_interface | default('eth0') }} \
    --rpc-port {{ node_rpc_port }} \
    --p2p-port {{ node_p2p_port }} \
    --metrics-port {{ node_metrics_port }} \
    --validator \
    --bootnodes {{ bootnodes | default('') }}
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ebpf-node

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ node_data_dir }} {{ node_log_dir }}
PrivateTmp=true

# Environment
Environment=RUST_LOG=info,libp2p=info,ebpf_node=info
Environment=NODE_ID={{ node_id | default('') }}
Environment=VALIDATOR_MODE={{ validator_mode | default('false') }}

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
