# services/ebpf-blockchain-node1.service
[Unit]
Description=eBPF Blockchain Node 1
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=maxi
Group=maxi
WorkingDirectory=/home/maxi/Documentos/source/ebpf-blockchain
ExecStart=/usr/bin/rustup run nightly /home/maxi/Documentos/source/ebpf-blockchain/ebpf-node/target/debug/ebpf-node \
    --iface eth0 \
    --listen-addresses /ip4/127.0.0.1/tcp/5001 \
    --bootstrap-peers /ip4/127.0.0.1/tcp/5002,/ip4/127.0.0.1/tcp/5003 \
    --data-dir /home/maxi/Documentos/source/ebpf-blockchain/data/node1 \
    --log-dir /home/maxi/Documentos/source/ebpf-blockchain/logs/node1 \
    --metrics-port 9090 \
    --rpc-port 9091 \
    --ws-port 9092
Restart=on-failure
RestartSec=10
StandardOutput=append:/home/maxi/Documentos/source/ebpf-blockchain/logs/node1/stdout.log
StandardError=append:/home/maxi/Documentos/source/ebpf-blockchain/logs/node1/stderr.log

# eBPF permissions
CapabilityBoundingSet=CAP_BPF CAP_PERFMON CAP_SYS_ADMIN
AmbientCapabilities=CAP_BPF CAP_PERFMON CAP_SYS_ADMIN
SecurityProtectHome=true

[Install]
WantedBy=multi-user.target
